# docker/supervisord.conf
; Supervisor 配置（管理 Nginx + PHP-FPM 进程）
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
logfile_maxbytes=10MB
loglevel=info
pidfile=/run/supervisord.pid

[unix_http_server]
file=/run/supervisor.sock
chmod=0770
chown=root:navwww

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///run/supervisor.sock

; ── dcron（计划任务，由 PHP 动态写入 navwww 的 crontab）──
[program:crond]
command=/usr/sbin/crond -f -l 8
autostart=true
autorestart=true
startretries=3
stdout_logfile=/var/log/crond.log
stderr_logfile=/var/log/crond.log
priority=20

; ── Nginx Reload Watcher（监听触发文件，执行 nginx reload）──
[program:nginx-reload-watcher]
command=/bin/sh -c 'while true; do if [ -f /tmp/nginx-reload-trigger ]; then rm -f /tmp/nginx-reload-trigger; kill -HUP $(cat /run/nginx/nginx.pid 2>/dev/null) 2>/dev/null; fi; sleep 1; done'
autostart=true
autorestart=true
startretries=99
stdout_logfile=/dev/null
stderr_logfile=/dev/null
priority=15

; ── Nginx ──
[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
startretries=3
stderr_logfile=/var/log/nginx/error.log
stdout_logfile=/var/log/nginx/access.log
priority=10

; ── PHP-FPM ──
[program:php-fpm]
command=/usr/local/sbin/php-fpm --nodaemonize --fpm-config /usr/local/etc/php-fpm.d/nav.conf
autostart=true
autorestart=true
startretries=3
stderr_logfile=/var/log/php-fpm/error.log
stdout_logfile=/var/log/php-fpm/error.log
priority=5
